# deploy/linux/hr-resume-filter.service
[Unit]
Description=HR Resume Filter
After=network.target

[Service]
Type=simple
User=hr-resume
Group=hr-resume
WorkingDirectory=/opt/hr-resume-filter
Environment=HR_RESUME_HOST=127.0.0.1
Environment=HR_RESUME_PORT=8765
ExecStart=/opt/hr-resume-filter/.venv/bin/python /opt/hr-resume-filter/web_app.py --no-browser
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
